# tasks/host/files/letsencrypt.service
[Unit]
Description=Certificate renewal

[Service]
Type=oneshot
ExecStart=/usr/bin/letsencrypt renew --webroot --webroot-path /var/www/letsencrypt
